# etc/profile.d/editor.sh
# Editor preference
if command -v nvim >/dev/null 2>&1; then
  export EDITOR="nvim"
elif command -v vim >/dev/null 2>&1; then
  export EDITOR="vim"
else
  export EDITOR="vi"
fi

# vim: set ft=sh ts=2 sw=2 et:
