# etc/profile.d/rust.sh
# Rust shell configuration
if [ -d ${HOME}/.cargo ]; then
  . "$HOME/.cargo/env"
fi

# vim: set ft=sh ts=2 sw=2 et:
